# extras/fclip/fclip.service
[Unit]
Description=Fuzzy Clipboard Watcher (Wayland)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/fclip watch
Restart=on-failure
RestartSec=30
StartLimitBurst=3

[Install]
WantedBy=default.target
